# cookbooks/asgard/templates/default/asgard.upstart.conf.erb
#!upstart
# my upstart asgard script
# this script will start/stop my asgard
# optional stuff
description "Asgard Server"

# configuration variables.
# You'll want to change thse as needed

start on filesystem and started networking

respawn
setuid asgard
env JAVA_HOME=<%= node['java']['java_home'] %>
exec java <%= node[:asgard][:java_options]%> -jar <%= node[:asgard][:base_dir]%>/asgard_standalone.jar 2>&1 >> <%= node[:asgard][:log_dir]%>/asgard.log
