# cookbooks/ice/templates/default/ice.upstart.conf.erb
#!upstart
# my upstart ICE script
# this script will start/stop my ICE
# optional stuff
description "ICE Server"

# configuration variables.
# You'll want to change thse as needed
env USER=ice
export JAVA_HOME=/usr/lib/jvm/default-java
export GRAILS_HOME=/usr/local/grails

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

respawn

chdir /usr/local/ice
exec sudo -u $USER sh -c "./grailsw run-app 2>&1 >> /var/log/ice/ice.log"
