# cookbooks/ice/templates/default/ice.upstart.conf.erb
#!upstart
# my upstart ICE script
# this script will start/stop my ICE
# optional stuff
description "ICE Server"

# configuration variables.
# You'll want to change thse as needed
env USER=ice

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

respawn

exec sudo -u $USER sh -c "java <%= node[:ice][:java_options]%> -jar <%= node[:ice][:base_dir]%>/asgard_standalone.jar 2>&1 >> <%= node[:ice][:log_dir]%>/ice.log"
