# cookbooks/ice/templates/default/ice.upstart.conf.erb
#!upstart
# my upstart ICE script
# this script will start/stop my ICE
# optional stuff
description "ICE Server"

# configuration variables.
# You'll want to change thse as needed

start on filesystem and started networking

respawn
setuid ice
env JAVA_HOME=<%= node['java']['java_home'] %>
env GRAILS_HOME=<%= node['grails']['home_dir'] %>
chdir /usr/local/ice
exec ./grailsw run-app 2>&1 >> <%= node['ice']['log_dir'] %>/ice.log
